# packaging/systemd/register-mvp.service
[Unit]
Description=Register MVP Service
After=network-online.target
Wants=network-online.target

[Service]
User=registermvp
Group=registermvp
ExecStart=/opt/register_mvp/bin/register_mvp --pos-http
AmbientCapabilities=
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallFilter=@system-service @network-io
ReadWritePaths=/var/log/register-mvp /var/lib/register-mvp /opt/register_mvp/data
RuntimeDirectory=register-mvp
StateDirectory=register-mvp
UMask=0077

[Install]
WantedBy=multi-user.target
